# scripts/systemd/octave-signal-engine.service
[Unit]
Description=Octave Signal Engine (3s detection loop)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=octave
WorkingDirectory=/home/octave/trading-alerts
ExecStart=/usr/bin/node --max-old-space-size=512 /home/octave/trading-alerts/src/index.js
Restart=always
RestartSec=10
EnvironmentFile=/home/octave/.config/trading-alerts/.env
StandardOutput=append:/home/octave/.octave-logs/signal-engine.log
StandardError=append:/home/octave/.octave-logs/signal-engine.log

[Install]
WantedBy=multi-user.target
